# supervisord.conf
[supervisord]
nodaemon=true
[program:shadowsocks-py]
command=ssserver -p 8388 -k nathan -m aes-256-cfb